# lesson 2.1.3/crontab/delete_data
0 23 * * * mysql -u lesson213 -plesson213 -h db_lesson213 -P 3306 -e "DELETE FROM lesson213.users"
30 23 * * * mysql -u lesson213 -plesson213 -h db_lesson213 -P 3306 -e "DELETE FROM lesson213.session"
