# lesson 1.3.2/crontab/delete_data
0 23 * * * mysql -u lesson132 -plesson132 -h db_lesson132 -P 3306 -e "DELETE FROM lesson132.users"
30 23 * * * mysql -u lesson132 -plesson132 -h db_lesson132 -P 3306 -e "DELETE FROM lesson132.session"
