# lesson 2.1.2/crontab/delete_data
0 23 * * * mysql -u lesson212 -plesson212 -h db_lesson212 -P 3306 -e "DELETE FROM lesson212.users"
30 23 * * * mysql -u lesson212 -plesson212 -h db_lesson212 -P 3306 -e "DELETE FROM lesson212.session"
