# lesson 1.3.1.3-4/crontab/delete_data
0 23 * * * mysql -u lesson13134 -plesson13134 -h db_lesson13134 -P 3306 -e "DELETE FROM lesson13134.users"
30 23 * * * mysql -u lesson13134 -plesson13134 -h db_lesson13134 -P 3306 -e "DELETE FROM lesson13134.session"
